# unix/dvorak@.service
[Unit]
Description=Dvorak Virtual Keyboard

[Service]
ExecStart=/usr/bin/dvorak /dev/input/%i -m "AT" 


[Install]
WantedBy=multi-user.target
